# supervisord.conf
[supervisord]
nodaemon=true

[program:tor]
command=/home/whoogle/whoogle/misc/tor/start-tor.sh

[program:whoogle]
command=/home/whoogle/whoogle/run
